# bitmediacentre-base-rpm/bitmediacentre-start.service
[Unit]

[Service]
Type=oneshot
ExecStartPre=chattr -i /
ExecStart=/bin/sh -c "[ -d /start ] || { ln -s /var/start /start; ln -s /var/start/env.sh /env.sh; mkdir -p /start/youtube; ln -s /usr/start/youtube/*.sh /var/start/ }"
ExecStopPost=chattr +i /

[Install]
WantedBy=basic.target
