# bitmediacentre-base-rpm/install-upgrade-youtube-dl.service
[Unit]
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/pip3 install --upgrade yt-dlp

[Install]
WantedBy=graphical.target
