# bitmediacentre-base-rpm/pulseaudio.service
[Unit]

[Service]
ExecStart=/usr/bin/pulseaudio --disallow-exit --exit-idle-time=-1
EnvironmentFile=/env.sh

[Install]
WantedBy=graphical.target
