# bitmediacentre-base-rpm/rc.service
[Unit]
After=network-online.target

[Service]
ExecStart=/bin/python3 /start/rc-server.py

[Install]
WantedBy=graphical.target
