# bitmediacentre-base-rpm/rc.service
[Unit]
Requires=myevemu-device.service

[Service]
ExecStart=/bin/python3 /usr/start/bin/rc-server2.py

[Install]
WantedBy=graphical.target
