# bitmediacentre-base-rpm/terminal.service
[Unit]
After=setup-gnome-terminal.service

[Service]
Type=oneshot
EnvironmentFile=/env.sh
ExecStart=/usr/bin/gnome-terminal --zoom=3.0 --hide-menubar -- /usr/start/bin/mymc
WorkingDirectory=/home/pcuser
