# bitmediacentre-base-rpm/terminal.timer
[Unit]
After=setup-gnome-terminal.service

[Timer]
OnActiveSec=8s
Unit=terminal.service

[Install]
WantedBy=graphical.target
