# bitmediacentre-base-rpm/terminal.timer
[Unit]
After=compositor.service

[Timer]
OnActiveSec=4s
Unit=terminal.service

[Install]
WantedBy=graphical.target
