# bitmediacentre-base-rpm/terminal.timer
[Unit]
#Note: not sure After= is working as intended here:
After=setup-gnome-terminal.service

[Timer]
OnActiveSec=10s
Unit=terminal.service

[Install]
WantedBy=graphical.target
