# bitmediacentre-base-rpm/terminal.service
[Unit]

[Service]
Type=simple
Restart=no
EnvironmentFile=/env.sh
ExecStart=/usr/start/bin/terminal.sh
WorkingDirectory=/home/pcuser
User=pcuser

